# scripts/com.example.hostshelper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.hostshelper</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/hosts-helper</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/var/log/hosts-helper.log</string>
    
    <key>StandardErrorPath</key>
    <string>/var/log/hosts-helper-error.log</string>
    
    <key>Sockets</key>
    <dict>
        <key>HostsHelperSocket</key>
        <dict>
            <key>SockPathName</key>
            <string>/var/run/com.example.hostshelper.sock</string>
            <key>SockPathMode</key>
            <integer>384</integer> <!-- 0600 permissions in decimal -->
        </dict>
    </dict>
    
    <key>UserName</key>
    <string>root</string>
    
    <key>GroupName</key>
    <string>wheel</string>
    
    <key>SessionCreate</key>
    <true/>
</dict>
</plist>